# ansible/roles/kafka/templates/kafka.service
[Unit]
Requires=zookeeper.service
After=zookeeper.service

[Service]
Type=simple
User=raj
ExecStart=/bin/sh -c '/home/raj/kafka/bin/kafka-server-start.sh /home/raj/kafka/config/server.properties > /home/raj/kafka/kafka.log 2>&1'
ExecStop=/home/raj/kafka/bin/kafka-server-stop.sh
Restart=on-abnormal

[Install]
WantedBy=multi-user.target