# ansible/roles/kafka/templates/zookeeper.service
[Unit]
Requires=network.target remote-fs.target
After=network.target remote-fs.target

[Service]
Type=simple
User=raj
ExecStart=/home/raj/kafka/bin/zookeeper-server-start.sh /home/raj/kafka/config/zookeeper.properties
ExecStop=/home/raj/kafka/bin/zookeeper-server-stop.sh
Restart=on-abnormal

[Install]
WantedBy=multi-user.target